# root/etc/services.d/tailscaled/finish
#!/usr/bin/with-contenv bash
# vi: ft=bash

tailscaled --cleanup
